# delay.timer
[Unit]
Description="Run script after 3 min"
[Timer]
OnBootSec=3m
[Install]
WantedBy=default.target

